# profile.d/extra.symlink
# Last things sourced save /etc/bashrc

# Fix CTRL-Q and CTRL-S in terminal
stty -ixon > /dev/null 2>&1

source ~/.profile.osx.d/extra
